# MQTTLogger.service
[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/matt/PycharmProjects/MQTTLogger/MQTTLogger.py

[Install]
WantedBy=multi-user.target
